# deploy/cast-bot.service
[Unit]
Description=CAST Telegram Bot
After=network.target postgresql.service

[Service]
User=cast
WorkingDirectory=/home/cast/cast
EnvironmentFile=/home/cast/cast/.env
ExecStart=/home/cast/cast/.venv/bin/python manage.py runbot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
